# templates/github-runner-auto-register.timer.j2
[Unit]
Description=GitHub Runner Auto-Register Timer

[Timer]
OnBootSec=1min
OnUnitActiveSec=15min

[Install]
WantedBy=timers.target 