# scripts/github-runner-auto-register.service
[Unit]
Description=GitHub Runner Auto-Registration Service
After=network.target
Wants=network-online.target

[Service]
Type=oneshot
User=root
Environment=GITHUB_TOKEN=%i
Environment=GITHUB_USERNAME=ivanzud
Environment=RUNNER_LABELS=self-hosted,linux,x64,ansible
Environment=RUNNER_NAME_PREFIX=auto-runner
Environment=RUNNER_WORK_DIR=/opt/github-runners
Environment=RUNNER_USER=github-runner
ExecStart=/opt/scripts/auto-register-runners.sh
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target 